# openvpn.target
[Unit]
Description=OpenVPN Networks

[Install]
WantedBy=multi-user.target
